# debian/postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable gpms-client-server-autoregi.service
