# debian/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" -a ! -d /run/systemd/system ]; then
    systemctl disable gpms-client-server-autoregi.service
fi


