# etc/systemd/system/gpms-client-server-autoregi.service
[Unit]
Description=Gooroom Platform client-server auto-registration
DefaultDependencies=gooroom-agent.service
After=gooroom-agent.service network-online.target 
Wants=network-online.target
Before=systemd-logind.service
ConditionPathExists=/usr/sbin/gpms-auto-registration.sh

[Service]
EnvironmentFile=/etc/gooroom/gooroom-client-server-register/gcsr_auto.conf 
ExecStart=/usr/sbin/gpms-auto-registration.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
